# supervisor-conf/develop/api/vero-api.ini
[program: api]
command=python manage.py runserver 0.0.0.0:8000
directory = /services/vero
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
priority=20

[program:migrations]
command=python manage.py migrate
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
directory = /services/vero
priority = 10
autorestart = false
